# deploy/dashboard-backup.timer
[Unit]
Description=Daily Datacenter Dashboard backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target