# deploy/dashboard.service
[Unit]
Description=Datacenter Dashboard API
Documentation=file:///opt/dashboard/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dashboard
Group=dashboard
WorkingDirectory=/opt/dashboard
EnvironmentFile=-/opt/dashboard/.env
Environment=PYTHONPATH=/opt/dashboard/backend
Environment=FRONTEND_STATIC_DIR=/opt/dashboard/frontend/dist
ExecStart=/opt/dashboard/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# SQLite DB, discovery probes, collector polling
NoNewPrivileges=true
ReadWritePaths=/opt/dashboard/data /opt/dashboard/backups

[Install]
WantedBy=multi-user.target