# deploy/dashboard.service.in
[Unit]
Description=Datacenter Dashboard API
Documentation=file://@INSTALL_DIR@/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@INSTALL_DIR@
EnvironmentFile=-@INSTALL_DIR@/.env
Environment=PYTHONPATH=@INSTALL_DIR@/backend
Environment=FRONTEND_STATIC_DIR=@INSTALL_DIR@/frontend/dist
ExecStart=@PYTHON@ -m uvicorn app.main:app --host @BIND_HOST@ --port @PORT@
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# SQLite DB, discovery probes, collector polling (ICMP via icmplib)
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
NoNewPrivileges=true
ReadWritePaths=@INSTALL_DIR@/data @INSTALL_DIR@/backups

[Install]
WantedBy=multi-user.target